# jarvis/scripts/g147.service
[Unit]
Description=g147

[Service]
ExecStart=/home/pepper/.g147/groot.sh

[Install]
WantedBy=multi-user.target
